# deploy/systemd/stcp-chatbot-discord.service
[Unit]
Description=STCP Chatbot Discord Gateway worker
Documentation=file:/opt/stcp-chatbot/docs/DISCORD.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stcp-chatbot
Group=stcp-chatbot
WorkingDirectory=/opt/stcp-chatbot
ExecStart=/usr/bin/php -d memory_limit=512M /opt/stcp-chatbot/bin/discord-worker.php
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/stcp-chatbot/storage

[Install]
WantedBy=multi-user.target
