# deploy/systemd/stcp-chatbot-matrix.service
[Unit]
Description=STCP Chatbot Matrix E2EE worker
Documentation=file:/opt/stcp-chatbot/docs/MATRIX.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stcp-chatbot
Group=stcp-chatbot
WorkingDirectory=/opt/stcp-chatbot
Environment=STCP_CHATBOT_PHP=/usr/bin/php
ExecStart=/usr/bin/node /opt/stcp-chatbot/matrix/worker.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/stcp-chatbot/storage

[Install]
WantedBy=multi-user.target
